# src/supervisord.conf
[supervisord]
nodaemon=true

[program:cassandra]
user=cassandra
command=cassandra -f
